# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:httpd]
command=service httpd start

[program:monast]
command=service monast start
