# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:httpd]
command=httpd

[program:monast]
command=monast
